# pi-tape/systemd/pi-tape-daemon.service
[Unit]
Description=Pi tape HTTP playback control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=audio
Environment=PI_TAPE_CONFIG=/etc/pi-tape/config.env
ExecStart=/usr/local/bin/pi-tape-daemon.py
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
